# packaging/systemd/toki-hub.service
[Unit]
Description=Toki encrypted snapshot hub
Documentation=https://github.com/choi138/toki/blob/main/docs/remote-sync.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=toki-hub
Group=toki-hub
StateDirectory=toki-hub
StateDirectoryMode=0700
RuntimeDirectory=toki-hub
RuntimeDirectoryMode=0750
EnvironmentFile=/etc/toki-hub.env
Environment=TOKI_HUB_STORAGE_PATH=/var/lib/toki-hub
Environment=TOKI_HUB_SOCKET_PATH=/run/toki-hub/toki-hub.sock
Environment=LOG_LEVEL=warning
ExecStart=/usr/local/bin/toki-hub
Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s
UMask=0007
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
MemoryMax=1G
TasksMax=256
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
