# packaging/systemd/toki-agent.service
[Unit]
Description=Toki usage collection agent
Documentation=https://github.com/choi138/toki/blob/main/docs/remote-sync.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStartPre=/usr/local/bin/toki-agent doctor
ExecStart=/usr/local/bin/toki-agent run
Restart=on-failure
RestartSec=30s
TimeoutStopSec=20s
UMask=0077
Nice=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
PrivateUsers=true
ProtectSystem=strict
ProtectHome=tmpfs
BindReadOnlyPaths=-%h/.claude/projects
BindReadOnlyPaths=-%h/.codex/state_5.sqlite
BindReadOnlyPaths=-%h/.codex/state_5.sqlite-wal
BindReadOnlyPaths=-%h/.codex/state_5.sqlite-shm
BindReadOnlyPaths=-%h/.codex/sessions
BindReadOnlyPaths=-%h/.codex/archived_sessions
BindReadOnlyPaths=-%h/.hermes/state.db
BindReadOnlyPaths=-%h/.hermes/state.db-wal
BindReadOnlyPaths=-%h/.hermes/state.db-shm
BindReadOnlyPaths=-%h/.config/Cursor/User/globalStorage/state.vscdb
BindReadOnlyPaths=-%h/.config/Cursor/User/globalStorage/state.vscdb-wal
BindReadOnlyPaths=-%h/.config/Cursor/User/globalStorage/state.vscdb-shm
BindReadOnlyPaths=-%h/.gemini/tmp
BindReadOnlyPaths=-%h/.gjc/agent/sessions
BindReadOnlyPaths=-%h/.local/share/opencode/opencode.db
BindReadOnlyPaths=-%h/.local/share/opencode/opencode.db-wal
BindReadOnlyPaths=-%h/.local/share/opencode/opencode.db-shm
BindReadOnlyPaths=-%h/.openclaw/agents
BindPaths=%h/.config/toki-agent
BindPaths=%h/.local/state/toki-agent
BindPaths=%h/.local/share/toki-agent
LockPersonality=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
MemoryMax=512M
TasksMax=128

[Install]
WantedBy=default.target
